# Python/PythonCache/gen-config/supervisor-devpi.conf

[program:devpi-server]
command=/usr/local/share/pip-global/bin/devpi-server --host 0.0.0.0
user = root
priority = 999
startsecs = 5
redirect_stderr = True
autostart = True
